# conf/etc/systemd/system/multi-user.target.wants/alexi-backpack.service
[Unit]
Description=Run Alexi's Backpack Access Point

[Service]
ExecStart=/usr/local/bin/alexi-backpack.sh start

[Install]
WantedBy=multi-user.target
